# supervisord.conf
; =============================================================================
; Supervisord — Crypto Trader single-container process manager
; Manages: MySQL, Redis, Celery worker, Celery beat, Flower, Flask
; =============================================================================

[supervisord]
nodaemon=true                       ; stay in foreground (required for Docker)
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── MySQL ─────────────────────────────────────────────────────────────────────
; Started by entrypoint.sh first (needs init); supervisor keeps it alive after.
[program:mysql]
command=/usr/sbin/mysqld --defaults-file=/etc/mysql/mysql.conf.d/zz-crypto.cnf
user=mysql
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=10                         ; start first
stdout_logfile=/data/logs/mysql-stdout.log
stderr_logfile=/data/logs/mysql-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Redis ─────────────────────────────────────────────────────────────────────
[program:redis]
command=redis-server /etc/redis/redis-crypto.conf
autostart=true
autorestart=true
startretries=5
startsecs=3
priority=20                         ; start after mysql
stdout_logfile=/data/logs/redis-stdout.log
stderr_logfile=/data/logs/redis-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Celery worker ─────────────────────────────────────────────────────────────
[program:celery-worker]
command=celery -A celery_app worker
    --loglevel=info
    --concurrency=4
    -Q data,analysis,trading
    --logfile=/data/logs/celery-worker.log
directory=/workspace
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=30                         ; start after MySQL + Redis are ready
stdout_logfile=/data/logs/celery-worker-stdout.log
stderr_logfile=/data/logs/celery-worker-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Celery beat (scheduler) ───────────────────────────────────────────────────
[program:celery-beat]
command=celery -A celery_app beat
    --loglevel=info
    --scheduler celery.beat:PersistentScheduler
    --schedule=/data/celerybeat-schedule
    --logfile=/data/logs/celery-beat.log
directory=/workspace
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=40                         ; start after worker is up
stdout_logfile=/data/logs/celery-beat-stdout.log
stderr_logfile=/data/logs/celery-beat-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Flower (Celery web monitor) ───────────────────────────────────────────────
[program:flower]
command=celery -A celery_app flower --port=5555 --logfile=/data/logs/flower.log
directory=/workspace
autostart=true
autorestart=true
startretries=3
startsecs=15
priority=50
stdout_logfile=/data/logs/flower-stdout.log
stderr_logfile=/data/logs/flower-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Flask dashboard ───────────────────────────────────────────────────────────
[program:flask]
command=flask --app app.dashboard.server run --host=0.0.0.0 --port=5000
directory=/workspace
autostart=true
autorestart=true
startretries=5
startsecs=20
priority=60                         ; start last — needs DB + Celery ready
stdout_logfile=/data/logs/flask-stdout.log
stderr_logfile=/data/logs/flask-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Group: restart all Python services together ───────────────────────────────
[group:python-services]
programs=celery-worker,celery-beat,flower,flask
priority=30
